# modules/ffmff/files/dns_server/reload-bind.path
[Unit]
Description=Watch for changes to the Freifunk config

[Path]
PathChanged=/var/lib/ffmff/output/zones.ff
