# modules/ffmff/files/dns_server/reload-bind.service
[Unit]
Description=Reload the bind Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload-or-try-restart bind9.service
