# modules/routing/templates/update-directexit.timer
[Unit]
Description=Trigger Population of directexit ipsets

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=multi-user.target
